# ansible/playbooks/osquery/templates/osqueryd.service.j2
[Unit]
Description=osqueryd
After=network.target

[Service]
Type=simple
Restart=on-failure
RestartSec=10
User=root
Group=root
EnvironmentFile={{ fleet_secret_env_file }}
ExecStop={{ osqueryd_path | regex_replace('\\/$', '') }}/osqueryctl stop
ExecStart={{ osqueryd_path | regex_replace('\\/$', '') }}/osqueryd \
  --flagfile={{ osqueryd_flags_file }}\
  --host_identifier={{ ansible_hostname }} \
  --tls_hostname={{ fleet_endpoint }} \
  --enroll_secret_env=ENROLL_SECRET

[Install]
WantedBy=multi-user.target
